# debian/postinst
#!/bin/bash
systemctl enable subutai-dnsmasq.service
systemctl start subutai-dnsmasq.service