# debian/preinst
#!/bin/bash
systemctl disable dnsmasq.service
systemctl stop dnsmasq.service